# roles/systemd_service/templates/unit.service.j2
[Unit]
Description={{ unit_description }}

[Service]
Type={{ service_type }}
ExecStart={{ service_exec_start }}

[Install]
WantedBy=multi-user.target
